# backend/deploy/systemd/coreops-backend.service
[Unit]
Description=CoreOps Backend (FastAPI / Uvicorn)
After=network.target postgresql.service redis.service
Wants=redis.service

[Service]
Type=simple
User=coreops
Group=coreops
WorkingDirectory=/opt/coreops/backend
EnvironmentFile=/opt/coreops/backend/.env
# Apply migrations on start, then run the API. The worker/beat units do NOT
# run migrations, so they run exactly once here.
ExecStartPre=/opt/coreops/backend/.venv/bin/alembic upgrade head
ExecStart=/opt/coreops/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8100
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
